# k8s-ansible-script/os/roles/ali-local-config/templates/alilocalconf.service
[Unit]
Description=collect system local config
After=network.target

[Service]
Type=oneshot
ExecStart={{ bin_dir }}/alilocalconf.sh

[Install]
WantedBy=multi-user.target