# k8s-ansible-script/k8s-cluster/roles/log-cleanup/files/kubelog-cleanup.timer
[Unit]
Description=clean up kubernets log files

[Timer]
OnCalendar=Mon *-*-* 00:00:00

[Install]
WantedBy=timers.target
